# telegram-bot.service
[Unit]
Description=Telegram Sticker Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/project
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
